# roles/nginx/templates/nginx.ini
[program:nginx]
command={{ user_yum_sbin_path }}/nginx -c {{ user_yum_etc_path }}/nginx/nginx.conf -p {{ user_yum_root_dir }}/usr/share/nginx
autostart=yes
autorestart=yes
stopasgroup=yes
killasgroup=yes
environment=LD_LIBRARY_PATH="{{ user_yum_root_dir }}/usr/lib:{{ user_yum_root_dir }}/usr/lib64"
