# roles/winnt4/templates/winnt4.ini
[program:winnt4]
directory={{ directory }}
command={{ user_bin }}/winnt4
autostart=yes
autorestart=yes
stopasgroup=yes
killasgroup=yes
